# usr/lib/systemd/system/rc-local.service
[Unit]
Description=Execute /etc/rc.local during startup

[Service]
Type=oneshot
ExecStart=/etc/rc.local

[Install]
WantedBy=multi-user.target
